# ext2fs/test3/var/lib/pacman/local/consolekit-0.4.6-4/files
%FILES%
etc/
etc/ConsoleKit/
etc/ConsoleKit/run-seat.d/
etc/ConsoleKit/run-session.d/
etc/ConsoleKit/seats.d/
etc/ConsoleKit/seats.d/00-primary.seat
etc/dbus-1/
etc/dbus-1/system.d/
etc/dbus-1/system.d/ConsoleKit.conf
etc/logrotate.d/
etc/logrotate.d/consolekit
usr/
usr/bin/
usr/bin/ck-history
usr/bin/ck-launch-session
usr/bin/ck-list-sessions
usr/include/
usr/include/ConsoleKit/
usr/include/ConsoleKit/ck-connector/
usr/include/ConsoleKit/ck-connector/ck-connector.h
usr/lib/
usr/lib/ConsoleKit/
usr/lib/ConsoleKit/ck-collect-session-info
usr/lib/ConsoleKit/ck-get-x11-display-device
usr/lib/ConsoleKit/ck-get-x11-server-pid
usr/lib/ConsoleKit/run-seat.d/
usr/lib/ConsoleKit/run-seat.d/udev-acl.ck
usr/lib/ConsoleKit/run-session.d/
usr/lib/ConsoleKit/run-session.d/pam-foreground-compat.ck
usr/lib/ConsoleKit/scripts/
usr/lib/ConsoleKit/scripts/ck-system-restart
usr/lib/ConsoleKit/scripts/ck-system-stop
usr/lib/libck-connector.so
usr/lib/libck-connector.so.0
usr/lib/libck-connector.so.0.0.0
usr/lib/pkgconfig/
usr/lib/pkgconfig/ck-connector.pc
usr/lib/security/
usr/lib/security/pam_ck_connector.so
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/basic.target.wants/
usr/lib/systemd/system/basic.target.wants/console-kit-log-system-start.service
usr/lib/systemd/system/console-kit-daemon.service
usr/lib/systemd/system/console-kit-log-system-restart.service
usr/lib/systemd/system/console-kit-log-system-start.service
usr/lib/systemd/system/console-kit-log-system-stop.service
usr/lib/systemd/system/halt.target.wants/
usr/lib/systemd/system/halt.target.wants/console-kit-log-system-stop.service
usr/lib/systemd/system/kexec.target.wants/
usr/lib/systemd/system/kexec.target.wants/console-kit-log-system-restart.service
usr/lib/systemd/system/poweroff.target.wants/
usr/lib/systemd/system/poweroff.target.wants/console-kit-log-system-stop.service
usr/lib/systemd/system/reboot.target.wants/
usr/lib/systemd/system/reboot.target.wants/console-kit-log-system-restart.service
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/consolekit.conf
usr/lib/udev/
usr/lib/udev/rules.d/
usr/lib/udev/rules.d/70-udev-acl.rules
usr/lib/udev/udev-acl
usr/sbin/
usr/sbin/ck-log-system-restart
usr/sbin/ck-log-system-start
usr/sbin/ck-log-system-stop
usr/sbin/console-kit-daemon
usr/share/
usr/share/dbus-1/
usr/share/dbus-1/interfaces/
usr/share/dbus-1/interfaces/org.freedesktop.ConsoleKit.Manager.xml
usr/share/dbus-1/interfaces/org.freedesktop.ConsoleKit.Seat.xml
usr/share/dbus-1/interfaces/org.freedesktop.ConsoleKit.Session.xml
usr/share/dbus-1/system-services/
usr/share/dbus-1/system-services/org.freedesktop.ConsoleKit.service
usr/share/doc/
usr/share/doc/ConsoleKit/
usr/share/doc/ConsoleKit/spec/
usr/share/doc/ConsoleKit/spec/ConsoleKit.html
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/pam_ck_connector.8.gz
usr/share/polkit-1/
usr/share/polkit-1/actions/
usr/share/polkit-1/actions/org.freedesktop.consolekit.policy
var/
var/log/
var/log/ConsoleKit/

