# partial_test2/var/lib/pacman/local/openssh-6.0p1-2/files
%FILES%
etc/
etc/conf.d/
etc/conf.d/sshd
etc/pam.d/
etc/pam.d/sshd
etc/rc.d/
etc/rc.d/functions.d/
etc/rc.d/functions.d/sshd-close-sessions
etc/rc.d/sshd
etc/ssh/
etc/ssh/moduli
etc/ssh/ssh_config
etc/ssh/sshd_config
usr/
usr/bin/
usr/bin/findssl.sh
usr/bin/scp
usr/bin/sftp
usr/bin/slogin
usr/bin/ssh
usr/bin/ssh-add
usr/bin/ssh-agent
usr/bin/ssh-copy-id
usr/bin/ssh-keygen
usr/bin/ssh-keyscan
usr/lib/
usr/lib/ssh/
usr/lib/ssh/sftp-server
usr/lib/ssh/ssh-keysign
usr/lib/ssh/ssh-pkcs11-helper
usr/sbin/
usr/sbin/sshd
usr/share/
usr/share/licenses/
usr/share/licenses/openssh/
usr/share/licenses/openssh/LICENCE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/scp.1.gz
usr/share/man/man1/sftp.1.gz
usr/share/man/man1/slogin.1.gz
usr/share/man/man1/ssh-add.1.gz
usr/share/man/man1/ssh-agent.1.gz
usr/share/man/man1/ssh-copy-id.1.gz
usr/share/man/man1/ssh-keygen.1.gz
usr/share/man/man1/ssh-keyscan.1.gz
usr/share/man/man1/ssh.1.gz
usr/share/man/man5/
usr/share/man/man5/moduli.5.gz
usr/share/man/man5/ssh_config.5.gz
usr/share/man/man5/sshd_config.5.gz
usr/share/man/man8/
usr/share/man/man8/sftp-server.8.gz
usr/share/man/man8/ssh-keysign.8.gz
usr/share/man/man8/ssh-pkcs11-helper.8.gz
usr/share/man/man8/sshd.8.gz
var/
var/empty/

%BACKUP%
etc/ssh/ssh_config	2c113824341eb4ab272c0b5b6bba1ae8
etc/ssh/sshd_config	387f21f67b7bb8cf1bdd6999006c1dd4
etc/pam.d/sshd	2fd20d311d2afbfc6e576883224d8c97
etc/conf.d/sshd	e2cea70ac13af7e63d40eb04415eacd5

