# partial_test2/var/lib/pacman/local/dbus-1.4.20-1/files
%FILES%
etc/
etc/X11/
etc/X11/xinit/
etc/X11/xinit/xinitrc.d/
etc/X11/xinit/xinitrc.d/30-dbus
usr/
usr/bin/
usr/bin/dbus-launch
usr/share/
usr/share/licenses/
usr/share/licenses/dbus/
usr/share/licenses/dbus/COPYING
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/dbus-launch.1.gz

