# ext2fs/test3/usr/lib/systemd/system/arch-daemons.target
[Unit]
Description=Arch Daemons
Documentation=man:arch-daemons(8)
