# ext2fs/var/lib/pacman/local/systemd-tools-185-1/files
%FILES%
etc/
etc/binfmt.d/
etc/modules-load.d/
etc/sysctl.d/
etc/tmpfiles.d/
etc/udev/
etc/udev/rules.d/
etc/udev/udev.conf
sbin/
sbin/udevadm
usr/
usr/bin/
usr/bin/systemd-ask-password
usr/bin/systemd-delta
usr/bin/systemd-detect-virt
usr/bin/systemd-tmpfiles
usr/bin/systemd-tty-ask-password-agent
usr/bin/udevadm
usr/bin/udevd
usr/include/
usr/include/gudev-1.0/
usr/include/gudev-1.0/gudev/
usr/include/gudev-1.0/gudev/gudev.h
usr/include/gudev-1.0/gudev/gudevclient.h
usr/include/gudev-1.0/gudev/gudevdevice.h
usr/include/gudev-1.0/gudev/gudevenumerator.h
usr/include/gudev-1.0/gudev/gudevenums.h
usr/include/gudev-1.0/gudev/gudevenumtypes.h
usr/include/gudev-1.0/gudev/gudevtypes.h
usr/include/libudev.h
usr/lib/
usr/lib/binfmt.d/
usr/lib/girepository-1.0/
usr/lib/girepository-1.0/GUdev-1.0.typelib
usr/lib/initcpio/
usr/lib/initcpio/hooks/
usr/lib/initcpio/hooks/udev
usr/lib/initcpio/install/
usr/lib/initcpio/install/timestamp
usr/lib/initcpio/install/udev
usr/lib/libgudev-1.0.a
usr/lib/libgudev-1.0.so
usr/lib/libgudev-1.0.so.0
usr/lib/libgudev-1.0.so.0.1.1
usr/lib/libudev.a
usr/lib/libudev.so
usr/lib/libudev.so.1
usr/lib/libudev.so.1.0.1
usr/lib/modules-load.d/
usr/lib/pkgconfig/
usr/lib/pkgconfig/gudev-1.0.pc
usr/lib/pkgconfig/libudev.pc
usr/lib/sysctl.d/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/sockets.target.wants/
usr/lib/systemd/system/sockets.target.wants/systemd-udev-control.socket
usr/lib/systemd/system/sockets.target.wants/systemd-udev-kernel.socket
usr/lib/systemd/system/sysinit.target.wants/
usr/lib/systemd/system/sysinit.target.wants/systemd-udev-trigger.service
usr/lib/systemd/system/sysinit.target.wants/systemd-udev.service
usr/lib/systemd/system/systemd-udev-control.socket
usr/lib/systemd/system/systemd-udev-kernel.socket
usr/lib/systemd/system/systemd-udev-settle.service
usr/lib/systemd/system/systemd-udev-trigger.service
usr/lib/systemd/system/systemd-udev.service
usr/lib/systemd/systemd-ac-power
usr/lib/systemd/systemd-binfmt
usr/lib/systemd/systemd-cryptsetup
usr/lib/systemd/systemd-modules-load
usr/lib/systemd/systemd-random-seed
usr/lib/systemd/systemd-remount-fs
usr/lib/systemd/systemd-reply-password
usr/lib/systemd/systemd-sysctl
usr/lib/systemd/systemd-timestamp
usr/lib/systemd/systemd-udevd
usr/lib/systemd/systemd-vconsole-setup
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/console.conf
usr/lib/tmpfiles.d/legacy.conf
usr/lib/tmpfiles.d/systemd.conf
usr/lib/tmpfiles.d/tmp.conf
usr/lib/tmpfiles.d/x11.conf
usr/lib/udev/
usr/lib/udev/accelerometer
usr/lib/udev/ata_id
usr/lib/udev/cdrom_id
usr/lib/udev/collect
usr/lib/udev/devices/
usr/lib/udev/findkeyboards
usr/lib/udev/keyboard-force-release.sh
usr/lib/udev/keymap
usr/lib/udev/keymaps/
usr/lib/udev/keymaps/acer
usr/lib/udev/keymaps/acer-aspire_5720
usr/lib/udev/keymaps/acer-aspire_5920g
usr/lib/udev/keymaps/acer-aspire_6920
usr/lib/udev/keymaps/acer-aspire_8930
usr/lib/udev/keymaps/acer-travelmate_c300
usr/lib/udev/keymaps/asus
usr/lib/udev/keymaps/compaq-e_evo
usr/lib/udev/keymaps/dell
usr/lib/udev/keymaps/dell-latitude-xt2
usr/lib/udev/keymaps/everex-xt5000
usr/lib/udev/keymaps/force-release/
usr/lib/udev/keymaps/force-release/common-volume-keys
usr/lib/udev/keymaps/force-release/dell-touchpad
usr/lib/udev/keymaps/force-release/dell-xps
usr/lib/udev/keymaps/force-release/hp-other
usr/lib/udev/keymaps/force-release/samsung-90x3a
usr/lib/udev/keymaps/force-release/samsung-other
usr/lib/udev/keymaps/fujitsu-amilo_li_2732
usr/lib/udev/keymaps/fujitsu-amilo_pa_2548
usr/lib/udev/keymaps/fujitsu-amilo_pro_edition_v3505
usr/lib/udev/keymaps/fujitsu-amilo_pro_v3205
usr/lib/udev/keymaps/fujitsu-amilo_si_1520
usr/lib/udev/keymaps/fujitsu-esprimo_mobile_v5
usr/lib/udev/keymaps/fujitsu-esprimo_mobile_v6
usr/lib/udev/keymaps/genius-slimstar-320
usr/lib/udev/keymaps/hewlett-packard
usr/lib/udev/keymaps/hewlett-packard-2510p_2530p
usr/lib/udev/keymaps/hewlett-packard-compaq_elitebook
usr/lib/udev/keymaps/hewlett-packard-pavilion
usr/lib/udev/keymaps/hewlett-packard-presario-2100
usr/lib/udev/keymaps/hewlett-packard-tablet
usr/lib/udev/keymaps/hewlett-packard-tx2
usr/lib/udev/keymaps/ibm-thinkpad-usb-keyboard-trackpoint
usr/lib/udev/keymaps/inventec-symphony_6.0_7.0
usr/lib/udev/keymaps/lenovo-3000
usr/lib/udev/keymaps/lenovo-ideapad
usr/lib/udev/keymaps/lenovo-thinkpad-usb-keyboard-trackpoint
usr/lib/udev/keymaps/lenovo-thinkpad_x200_tablet
usr/lib/udev/keymaps/lenovo-thinkpad_x6_tablet
usr/lib/udev/keymaps/lg-x110
usr/lib/udev/keymaps/logitech-wave
usr/lib/udev/keymaps/logitech-wave-cordless
usr/lib/udev/keymaps/logitech-wave-pro-cordless
usr/lib/udev/keymaps/maxdata-pro_7000
usr/lib/udev/keymaps/medion-fid2060
usr/lib/udev/keymaps/medionnb-a555
usr/lib/udev/keymaps/micro-star
usr/lib/udev/keymaps/module-asus-w3j
usr/lib/udev/keymaps/module-ibm
usr/lib/udev/keymaps/module-lenovo
usr/lib/udev/keymaps/module-sony
usr/lib/udev/keymaps/module-sony-old
usr/lib/udev/keymaps/module-sony-vgn
usr/lib/udev/keymaps/olpc-xo
usr/lib/udev/keymaps/onkyo
usr/lib/udev/keymaps/oqo-model2
usr/lib/udev/keymaps/samsung-90x3a
usr/lib/udev/keymaps/samsung-other
usr/lib/udev/keymaps/samsung-sq1us
usr/lib/udev/keymaps/samsung-sx20s
usr/lib/udev/keymaps/toshiba-satellite_a100
usr/lib/udev/keymaps/toshiba-satellite_a110
usr/lib/udev/keymaps/toshiba-satellite_m30x
usr/lib/udev/keymaps/zepto-znote
usr/lib/udev/mtd_probe
usr/lib/udev/rules.d/
usr/lib/udev/rules.d/42-usb-hid-pm.rules
usr/lib/udev/rules.d/50-udev-default.rules
usr/lib/udev/rules.d/60-cdrom_id.rules
usr/lib/udev/rules.d/60-persistent-alsa.rules
usr/lib/udev/rules.d/60-persistent-input.rules
usr/lib/udev/rules.d/60-persistent-serial.rules
usr/lib/udev/rules.d/60-persistent-storage-tape.rules
usr/lib/udev/rules.d/60-persistent-storage.rules
usr/lib/udev/rules.d/60-persistent-v4l.rules
usr/lib/udev/rules.d/61-accelerometer.rules
usr/lib/udev/rules.d/70-power-switch.rules
usr/lib/udev/rules.d/70-uaccess.rules
usr/lib/udev/rules.d/71-seat.rules
usr/lib/udev/rules.d/73-seat-late.rules
usr/lib/udev/rules.d/75-net-description.rules
usr/lib/udev/rules.d/75-probe_mtd.rules
usr/lib/udev/rules.d/75-tty-description.rules
usr/lib/udev/rules.d/78-sound-card.rules
usr/lib/udev/rules.d/80-drivers.rules
usr/lib/udev/rules.d/95-keyboard-force-release.rules
usr/lib/udev/rules.d/95-keymap.rules
usr/lib/udev/rules.d/95-udev-late.rules
usr/lib/udev/rules.d/99-systemd.rules
usr/lib/udev/scsi_id
usr/lib/udev/udevd
usr/lib/udev/v4l_id
usr/share/
usr/share/gir-1.0/
usr/share/gir-1.0/GUdev-1.0.gir
usr/share/gtk-doc/
usr/share/gtk-doc/html/
usr/share/gtk-doc/html/gudev/
usr/share/gtk-doc/html/gudev/GUdevClient.html
usr/share/gtk-doc/html/gudev/GUdevDevice.html
usr/share/gtk-doc/html/gudev/GUdevEnumerator.html
usr/share/gtk-doc/html/gudev/annotation-glossary.html
usr/share/gtk-doc/html/gudev/api-index-deprecated.html
usr/share/gtk-doc/html/gudev/api-index-full.html
usr/share/gtk-doc/html/gudev/gudev-hierarchy.html
usr/share/gtk-doc/html/gudev/gudev.devhelp2
usr/share/gtk-doc/html/gudev/home.png
usr/share/gtk-doc/html/gudev/index.html
usr/share/gtk-doc/html/gudev/index.sgml
usr/share/gtk-doc/html/gudev/ix02.html
usr/share/gtk-doc/html/gudev/left.png
usr/share/gtk-doc/html/gudev/ref-API.html
usr/share/gtk-doc/html/gudev/right.png
usr/share/gtk-doc/html/gudev/style.css
usr/share/gtk-doc/html/gudev/up.png
usr/share/gtk-doc/html/libudev/
usr/share/gtk-doc/html/libudev/api-index-full.html
usr/share/gtk-doc/html/libudev/ch01.html
usr/share/gtk-doc/html/libudev/home.png
usr/share/gtk-doc/html/libudev/index.html
usr/share/gtk-doc/html/libudev/index.sgml
usr/share/gtk-doc/html/libudev/left.png
usr/share/gtk-doc/html/libudev/libudev-udev-device.html
usr/share/gtk-doc/html/libudev/libudev-udev-enumerate.html
usr/share/gtk-doc/html/libudev/libudev-udev-list.html
usr/share/gtk-doc/html/libudev/libudev-udev-monitor.html
usr/share/gtk-doc/html/libudev/libudev-udev-queue.html
usr/share/gtk-doc/html/libudev/libudev-udev-util.html
usr/share/gtk-doc/html/libudev/libudev-udev.html
usr/share/gtk-doc/html/libudev/libudev.devhelp2
usr/share/gtk-doc/html/libudev/right.png
usr/share/gtk-doc/html/libudev/style.css
usr/share/gtk-doc/html/libudev/up.png
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/systemd-ask-password.1.gz
usr/share/man/man1/systemd-delta.1.gz
usr/share/man/man1/systemd-detect-virt.1.gz
usr/share/man/man5/
usr/share/man/man5/binfmt.d.5.gz
usr/share/man/man5/hostname.5.gz
usr/share/man/man5/locale.conf.5.gz
usr/share/man/man5/modules-load.d.5.gz
usr/share/man/man5/sysctl.d.5.gz
usr/share/man/man5/tmpfiles.d.5.gz
usr/share/man/man5/vconsole.conf.5.gz
usr/share/man/man7/
usr/share/man/man7/udev.7.gz
usr/share/man/man8/
usr/share/man/man8/systemd-udevd.8.gz
usr/share/man/man8/udevadm.8.gz
usr/share/pkgconfig/
usr/share/pkgconfig/udev.pc

