# partial_test2/var/lib/pacman/local/cronie-1.4.8-1/files
%FILES%
etc/
etc/anacrontab
etc/conf.d/
etc/conf.d/crond
etc/cron.d/
etc/cron.d/0hourly
etc/cron.daily/
etc/cron.deny
etc/cron.hourly/
etc/cron.hourly/0anacron
etc/cron.monthly/
etc/cron.weekly/
etc/pam.d/
etc/pam.d/crond
etc/rc.d/
etc/rc.d/crond
usr/
usr/bin/
usr/bin/crontab
usr/sbin/
usr/sbin/anacron
usr/sbin/crond
usr/share/
usr/share/licenses/
usr/share/licenses/cronie/
usr/share/licenses/cronie/COPYING
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/crontab.1.gz
usr/share/man/man5/
usr/share/man/man5/anacrontab.5.gz
usr/share/man/man5/crontab.5.gz
usr/share/man/man8/
usr/share/man/man8/anacron.8.gz
usr/share/man/man8/cron.8.gz
usr/share/man/man8/crond.8.gz
var/
var/spool/
var/spool/anacron/
var/spool/cron/

%BACKUP%
etc/anacrontab	8241db83d5edf01c71734e41e383e205
etc/conf.d/crond	bf644ac15632af6a3788e5de7be13080
etc/pam.d/crond	ace2584e5467371e35acf5f5af8f3c29
etc/cron.deny	2962d9812fbbbdee1bd7f1c1ccb5e372

