# partial_test2/var/lib/pacman/local/xmms2-0.8DrO_o-6/files
%FILES%
etc/
etc/conf.d/
etc/conf.d/xmms2d.conf
etc/rc.d/
etc/rc.d/xmms2d
usr/
usr/bin/
usr/bin/nyxmms2
usr/bin/xmms2
usr/bin/xmms2-et
usr/bin/xmms2-find-avahi
usr/bin/xmms2-launcher
usr/bin/xmms2-mdns-avahi
usr/bin/xmms2-mlib-updater
usr/bin/xmms2d
usr/include/
usr/include/xmms2/
usr/include/xmms2/xmms/
usr/include/xmms2/xmms/xmms_bindata.h
usr/include/xmms2/xmms/xmms_config.h
usr/include/xmms2/xmms/xmms_error.h
usr/include/xmms2/xmms/xmms_ipc.h
usr/include/xmms2/xmms/xmms_log.h
usr/include/xmms2/xmms/xmms_medialib.h
usr/include/xmms2/xmms/xmms_object.h
usr/include/xmms2/xmms/xmms_outputplugin.h
usr/include/xmms2/xmms/xmms_plugin.h
usr/include/xmms2/xmms/xmms_sample.h
usr/include/xmms2/xmms/xmms_streamtype.h
usr/include/xmms2/xmms/xmms_util.h
usr/include/xmms2/xmms/xmms_xformplugin.h
usr/include/xmms2/xmmsc/
usr/include/xmms2/xmmsc/xmmsc_compiler.h
usr/include/xmms2/xmmsc/xmmsc_errorcodes.h
usr/include/xmms2/xmmsc/xmmsc_idnumbers.h
usr/include/xmms2/xmmsc/xmmsc_ipc_msg.h
usr/include/xmms2/xmmsc/xmmsc_ipc_transport.h
usr/include/xmms2/xmmsc/xmmsc_sockets.h
usr/include/xmms2/xmmsc/xmmsc_stdbool.h
usr/include/xmms2/xmmsc/xmmsc_stdint.h
usr/include/xmms2/xmmsc/xmmsc_stringport.h
usr/include/xmms2/xmmsc/xmmsc_strlist.h
usr/include/xmms2/xmmsc/xmmsc_unistd.h
usr/include/xmms2/xmmsc/xmmsc_util.h
usr/include/xmms2/xmmsc/xmmsc_visualization.h
usr/include/xmms2/xmmsc/xmmsv.h
usr/include/xmms2/xmmsc/xmmsv_bitbuffer.h
usr/include/xmms2/xmmsc/xmmsv_build.h
usr/include/xmms2/xmmsc/xmmsv_coll.h
usr/include/xmms2/xmmsc/xmmsv_deprecated.h
usr/include/xmms2/xmmsc/xmmsv_dict.h
usr/include/xmms2/xmmsc/xmmsv_general.h
usr/include/xmms2/xmmsc/xmmsv_list.h
usr/include/xmms2/xmmsc/xmmsv_util.h
usr/include/xmms2/xmmsclient/
usr/include/xmms2/xmmsclient/xmmsclient++-glib.h
usr/include/xmms2/xmmsclient/xmmsclient++.h
usr/include/xmms2/xmmsclient/xmmsclient++/
usr/include/xmms2/xmmsclient/xmmsclient++/bindata.h
usr/include/xmms2/xmmsclient/xmmsclient++/client.h
usr/include/xmms2/xmmsclient/xmmsclient++/coll.h
usr/include/xmms2/xmmsclient/xmmsclient++/collection.h
usr/include/xmms2/xmmsclient/xmmsclient++/config.h
usr/include/xmms2/xmmsclient/xmmsclient++/dict.h
usr/include/xmms2/xmmsclient/xmmsclient++/exceptions.h
usr/include/xmms2/xmmsclient/xmmsclient++/helpers.h
usr/include/xmms2/xmmsclient/xmmsclient++/list.h
usr/include/xmms2/xmmsclient/xmmsclient++/listener.h
usr/include/xmms2/xmmsclient/xmmsclient++/mainloop.h
usr/include/xmms2/xmmsclient/xmmsclient++/medialib.h
usr/include/xmms2/xmmsclient/xmmsclient++/playback.h
usr/include/xmms2/xmmsclient/xmmsclient++/playlist.h
usr/include/xmms2/xmmsclient/xmmsclient++/result.h
usr/include/xmms2/xmmsclient/xmmsclient++/signal.h
usr/include/xmms2/xmmsclient/xmmsclient++/stats.h
usr/include/xmms2/xmmsclient/xmmsclient++/typedefs.h
usr/include/xmms2/xmmsclient/xmmsclient++/xform.h
usr/include/xmms2/xmmsclient/xmmsclient-cf.h
usr/include/xmms2/xmmsclient/xmmsclient-ecore.h
usr/include/xmms2/xmmsclient/xmmsclient-glib.h
usr/include/xmms2/xmmsclient/xmmsclient-qt.h
usr/include/xmms2/xmmsclient/xmmsclient.h
usr/lib/
usr/lib/libxmmsclient++-glib.so
usr/lib/libxmmsclient++-glib.so.1
usr/lib/libxmmsclient++-glib.so.1.0.0
usr/lib/libxmmsclient++.so
usr/lib/libxmmsclient++.so.4
usr/lib/libxmmsclient++.so.4.0.0
usr/lib/libxmmsclient-glib.so
usr/lib/libxmmsclient-glib.so.1
usr/lib/libxmmsclient-glib.so.1.0.0
usr/lib/libxmmsclient.so
usr/lib/libxmmsclient.so.6
usr/lib/libxmmsclient.so.6.0.0
usr/lib/perl5/
usr/lib/perl5/site_perl/
usr/lib/perl5/site_perl/Audio/
usr/lib/perl5/site_perl/Audio/XMMSClient.pm
usr/lib/perl5/site_perl/Audio/XMMSClient/
usr/lib/perl5/site_perl/Audio/XMMSClient/Collection.pm
usr/lib/perl5/site_perl/Audio/XMMSClient/Glib.pm
usr/lib/perl5/site_perl/Audio/XMMSClient/Sync.pm
usr/lib/perl5/site_perl/auto/
usr/lib/perl5/site_perl/auto/Audio/
usr/lib/perl5/site_perl/auto/Audio/XMMSClient/
usr/lib/perl5/site_perl/auto/Audio/XMMSClient/XMMSClient.so
usr/lib/pkgconfig/
usr/lib/pkgconfig/xmms2-client-cpp-glib.pc
usr/lib/pkgconfig/xmms2-client-cpp.pc
usr/lib/pkgconfig/xmms2-client-glib.pc
usr/lib/pkgconfig/xmms2-client.pc
usr/lib/pkgconfig/xmms2-plugin.pc
usr/lib/python2.7/
usr/lib/python2.7/site-packages/
usr/lib/python2.7/site-packages/xmmsclient/
usr/lib/python2.7/site-packages/xmmsclient/__init__.py
usr/lib/python2.7/site-packages/xmmsclient/collections.py
usr/lib/python2.7/site-packages/xmmsclient/consts.py
usr/lib/python2.7/site-packages/xmmsclient/glib.py
usr/lib/python2.7/site-packages/xmmsclient/propdict.py
usr/lib/python2.7/site-packages/xmmsclient/qt3.py
usr/lib/python2.7/site-packages/xmmsclient/sync.py
usr/lib/python2.7/site-packages/xmmsclient/xmmsapi.so
usr/lib/python2.7/site-packages/xmmsclient/xmmsvalue.so
usr/lib/python3.2/
usr/lib/python3.2/site-packages/
usr/lib/python3.2/site-packages/xmmsclient/
usr/lib/python3.2/site-packages/xmmsclient/__init__.py
usr/lib/python3.2/site-packages/xmmsclient/collections.py
usr/lib/python3.2/site-packages/xmmsclient/consts.py
usr/lib/python3.2/site-packages/xmmsclient/glib.py
usr/lib/python3.2/site-packages/xmmsclient/propdict.py
usr/lib/python3.2/site-packages/xmmsclient/qt3.py
usr/lib/python3.2/site-packages/xmmsclient/sync.py
usr/lib/python3.2/site-packages/xmmsclient/xmmsapi.cpython-32mu.so
usr/lib/python3.2/site-packages/xmmsclient/xmmsvalue.cpython-32mu.so
usr/lib/ruby/
usr/lib/ruby/vendor_ruby/
usr/lib/ruby/vendor_ruby/1.9.1/
usr/lib/ruby/vendor_ruby/1.9.1/armv5tel-linux-eabi/
usr/lib/ruby/vendor_ruby/1.9.1/armv5tel-linux-eabi/xmmsclient_ext.so
usr/lib/ruby/vendor_ruby/1.9.1/armv5tel-linux-eabi/xmmsclient_glib.so
usr/lib/ruby/vendor_ruby/1.9.1/xmmsclient.rb
usr/lib/ruby/vendor_ruby/1.9.1/xmmsclient/
usr/lib/ruby/vendor_ruby/1.9.1/xmmsclient/async.rb
usr/lib/ruby/vendor_ruby/1.9.1/xmmsclient/sync.rb
usr/lib/xmms2/
usr/lib/xmms2/libxmms_airplay.so
usr/lib/xmms2/libxmms_alsa.so
usr/lib/xmms2/libxmms_ao.so
usr/lib/xmms2/libxmms_apefile.so
usr/lib/xmms2/libxmms_asf.so
usr/lib/xmms2/libxmms_asx.so
usr/lib/xmms2/libxmms_avcodec.so
usr/lib/xmms2/libxmms_cue.so
usr/lib/xmms2/libxmms_curl.so
usr/lib/xmms2/libxmms_daap.so
usr/lib/xmms2/libxmms_diskwrite.so
usr/lib/xmms2/libxmms_equalizer.so
usr/lib/xmms2/libxmms_faad.so
usr/lib/xmms2/libxmms_file.so
usr/lib/xmms2/libxmms_flac.so
usr/lib/xmms2/libxmms_flv.so
usr/lib/xmms2/libxmms_gme.so
usr/lib/xmms2/libxmms_gvfs.so
usr/lib/xmms2/libxmms_html.so
usr/lib/xmms2/libxmms_ices.so
usr/lib/xmms2/libxmms_icymetaint.so
usr/lib/xmms2/libxmms_id3v2.so
usr/lib/xmms2/libxmms_jack.so
usr/lib/xmms2/libxmms_karaoke.so
usr/lib/xmms2/libxmms_m3u.so
usr/lib/xmms2/libxmms_mad.so
usr/lib/xmms2/libxmms_mms.so
usr/lib/xmms2/libxmms_modplug.so
usr/lib/xmms2/libxmms_mp4.so
usr/lib/xmms2/libxmms_mpg123.so
usr/lib/xmms2/libxmms_musepack.so
usr/lib/xmms2/libxmms_normalize.so
usr/lib/xmms2/libxmms_null.so
usr/lib/xmms2/libxmms_nulstripper.so
usr/lib/xmms2/libxmms_ofa.so
usr/lib/xmms2/libxmms_oss.so
usr/lib/xmms2/libxmms_pls.so
usr/lib/xmms2/libxmms_pulse.so
usr/lib/xmms2/libxmms_replaygain.so
usr/lib/xmms2/libxmms_rss.so
usr/lib/xmms2/libxmms_samba.so
usr/lib/xmms2/libxmms_sndfile.so
usr/lib/xmms2/libxmms_speex.so
usr/lib/xmms2/libxmms_tta.so
usr/lib/xmms2/libxmms_vocoder.so
usr/lib/xmms2/libxmms_vorbis.so
usr/lib/xmms2/libxmms_wave.so
usr/lib/xmms2/libxmms_wavpack.so
usr/lib/xmms2/libxmms_xml.so
usr/lib/xmms2/libxmms_xspf.so
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/xmms2-et.1.gz
usr/share/man/man1/xmms2-launcher.1.gz
usr/share/man/man1/xmms2-mdns-avahi.1.gz
usr/share/man/man1/xmms2.1.gz
usr/share/man/man1/xmms2d.1.gz
usr/share/pixmaps/
usr/share/pixmaps/xmms2-128.png
usr/share/pixmaps/xmms2-16.png
usr/share/pixmaps/xmms2-32.png
usr/share/pixmaps/xmms2-48.png
usr/share/pixmaps/xmms2-black-on-white.svg
usr/share/pixmaps/xmms2-white-on-black.svg
usr/share/pixmaps/xmms2.svg
usr/share/xmms2/
usr/share/xmms2/mind.in.a.box-lament_snipplet.ogg
usr/share/xmms2/scripts/
usr/share/xmms2/scripts/startup.d/
usr/share/xmms2/scripts/startup.d/xmms2-et-launcher.sh
usr/share/xmms2/scripts/startup.d/xmms2-mdns-launcher.sh
usr/share/xmms2/scripts/startup.d/xmms2-mlib-updater.sh

%BACKUP%
etc/conf.d/xmms2d.conf	af13c937bf3c86b77ae6820107aab9b8

