# ext2fs/var/lib/pacman/local/initscripts-2012.10.1-1/files
%FILES%
etc/
etc/inittab
etc/logrotate.d/
etc/logrotate.d/bootlog
etc/rc.conf
etc/rc.d/
etc/rc.d/functions
etc/rc.d/functions.d/
etc/rc.d/hwclock
etc/rc.d/netfs
etc/rc.d/network
etc/rc.local
etc/rc.local.shutdown
etc/rc.multi
etc/rc.shutdown
etc/rc.single
etc/rc.sysinit
usr/
usr/lib/
usr/lib/systemd/
usr/lib/systemd/arch-modules-load
usr/lib/systemd/system-generators/
usr/lib/systemd/system-generators/arch-daemons
usr/lib/systemd/system/
usr/lib/systemd/system/arch-daemons.target
usr/lib/systemd/system/multi-user.target.wants/
usr/lib/systemd/system/multi-user.target.wants/arch-daemons.target
usr/lib/systemd/system/multi-user.target.wants/rc-local.service
usr/lib/systemd/system/netfs.service
usr/lib/systemd/system/rc-local-shutdown.service
usr/lib/systemd/system/rc-local.service
usr/lib/systemd/system/shutdown.target.wants/
usr/lib/systemd/system/shutdown.target.wants/rc-local-shutdown.service
usr/lib/systemd/system/sysinit.target.wants/
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/initscripts.conf
usr/sbin/
usr/sbin/rc.d
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/rc.d
usr/share/man/
usr/share/man/man5/
usr/share/man/man5/rc.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/arch-daemons.8.gz
usr/share/man/man8/arch-modules-load.8.gz
usr/share/man/man8/rc.d.8.gz
usr/share/zsh/
usr/share/zsh/site-functions/
usr/share/zsh/site-functions/_rc.d

%BACKUP%
etc/inittab	f02566fd9a1eeac3aadd5d6af32453c5
etc/rc.conf	c618e099f522f861e08a475f8960a6c2
etc/rc.local	617b7829ad58893f3c044bbc22e9e142
etc/rc.local.shutdown	669ce6d766dd9deb86686219f3763fe4

