# partial_test2/var/lib/pacman/local/initscripts-2012.06.1-1/files
%FILES%
etc/
etc/conf.d/
etc/conf.d/wireless
etc/inittab
etc/logrotate.d/
etc/logrotate.d/bootlog
etc/profile.d/
etc/profile.d/locale.sh
etc/rc.conf
etc/rc.d/
etc/rc.d/functions
etc/rc.d/functions.d/
etc/rc.d/hwclock
etc/rc.d/netfs
etc/rc.d/network
etc/rc.local
etc/rc.local.shutdown
etc/rc.multi
etc/rc.shutdown
etc/rc.single
etc/rc.sysinit
usr/
usr/lib/
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/initscripts.conf
usr/sbin/
usr/sbin/rc.d
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/rc.d
usr/share/man/
usr/share/man/man5/
usr/share/man/man5/rc.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/rc.d.8.gz
usr/share/zsh/
usr/share/zsh/site-functions/
usr/share/zsh/site-functions/_rc.d

%BACKUP%
etc/inittab	f02566fd9a1eeac3aadd5d6af32453c5
etc/rc.conf	c618e099f522f861e08a475f8960a6c2
etc/rc.local	bd236a57bc00aa4727b98c5ae401c457
etc/rc.local.shutdown	669ce6d766dd9deb86686219f3763fe4
etc/conf.d/wireless	027576534885b8d5dded9be546057b12

