# ext2fs/var/lib/pacman/local/syslog-ng-3.3.5-1/files
%FILES%
etc/
etc/conf.d/
etc/conf.d/syslog-ng
etc/logrotate.d/
etc/logrotate.d/syslog-ng
etc/rc.d/
etc/rc.d/syslog-ng
etc/syslog-ng/
etc/syslog-ng/modules.conf
etc/syslog-ng/patterndb.d/
etc/syslog-ng/scl.conf
etc/syslog-ng/syslog-ng.conf
usr/
usr/bin/
usr/bin/loggen
usr/bin/pdbtool
usr/bin/update-patterndb
usr/include/
usr/include/syslog-ng/
usr/include/syslog-ng/afinter.h
usr/include/syslog-ng/alarms.h
usr/include/syslog-ng/apphook.h
usr/include/syslog-ng/atomic.h
usr/include/syslog-ng/block-ref-parser.h
usr/include/syslog-ng/center.h
usr/include/syslog-ng/cfg-grammar.h
usr/include/syslog-ng/cfg-lexer.h
usr/include/syslog-ng/cfg-parser.h
usr/include/syslog-ng/cfg.h
usr/include/syslog-ng/children.h
usr/include/syslog-ng/compat.h
usr/include/syslog-ng/control.h
usr/include/syslog-ng/crypto.h
usr/include/syslog-ng/dgroup.h
usr/include/syslog-ng/dnscache.h
usr/include/syslog-ng/driver.h
usr/include/syslog-ng/filter-expr-parser.h
usr/include/syslog-ng/filter.h
usr/include/syslog-ng/gprocess.h
usr/include/syslog-ng/gsockaddr.h
usr/include/syslog-ng/gsocket.h
usr/include/syslog-ng/logmatcher.h
usr/include/syslog-ng/logmpx.h
usr/include/syslog-ng/logmsg.h
usr/include/syslog-ng/logparser.h
usr/include/syslog-ng/logpipe.h
usr/include/syslog-ng/logprocess.h
usr/include/syslog-ng/logproto.h
usr/include/syslog-ng/logqueue-fifo.h
usr/include/syslog-ng/logqueue.h
usr/include/syslog-ng/logreader.h
usr/include/syslog-ng/logrewrite.h
usr/include/syslog-ng/logsource.h
usr/include/syslog-ng/logstamp.h
usr/include/syslog-ng/logtransport.h
usr/include/syslog-ng/logwriter.h
usr/include/syslog-ng/mainloop.h
usr/include/syslog-ng/memtrace.h
usr/include/syslog-ng/messages.h
usr/include/syslog-ng/misc.h
usr/include/syslog-ng/msg-format.h
usr/include/syslog-ng/nvtable.h
usr/include/syslog-ng/parser-expr-parser.h
usr/include/syslog-ng/persist-state.h
usr/include/syslog-ng/plugin.h
usr/include/syslog-ng/pragma-parser.h
usr/include/syslog-ng/rewrite-expr-parser.h
usr/include/syslog-ng/scratch-buffers.h
usr/include/syslog-ng/serialize.h
usr/include/syslog-ng/sgroup.h
usr/include/syslog-ng/stats.h
usr/include/syslog-ng/str-format.h
usr/include/syslog-ng/syslog-names.h
usr/include/syslog-ng/syslog-ng.h
usr/include/syslog-ng/tags.h
usr/include/syslog-ng/templates.h
usr/include/syslog-ng/timeutils.h
usr/include/syslog-ng/tls-support.h
usr/include/syslog-ng/tlscontext.h
usr/include/syslog-ng/tlstransport.h
usr/include/syslog-ng/utils.h
usr/include/syslog-ng/value-pairs.h
usr/lib/
usr/lib/libsyslog-ng-3.3.5.so
usr/lib/libsyslog-ng.so
usr/lib/pkgconfig/
usr/lib/pkgconfig/syslog-ng.pc
usr/lib/syslog-ng/
usr/lib/syslog-ng/libaffile.so
usr/lib/syslog-ng/libafmongodb.so
usr/lib/syslog-ng/libafprog.so
usr/lib/syslog-ng/libafsocket-notls.so
usr/lib/syslog-ng/libafsocket-tls.so
usr/lib/syslog-ng/libafsocket.so
usr/lib/syslog-ng/libafuser.so
usr/lib/syslog-ng/libbasicfuncs.so
usr/lib/syslog-ng/libconfgen.so
usr/lib/syslog-ng/libconvertfuncs.so
usr/lib/syslog-ng/libcsvparser.so
usr/lib/syslog-ng/libdbparser.so
usr/lib/syslog-ng/libdummy.so
usr/lib/syslog-ng/libsyslog-ng-crypto.so
usr/lib/syslog-ng/libsyslogformat.so
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/syslog-ng.service
usr/sbin/
usr/sbin/syslog-ng
usr/sbin/syslog-ng-ctl
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/loggen.1.gz
usr/share/man/man1/pdbtool.1.gz
usr/share/man/man1/syslog-ng-ctl.1.gz
usr/share/man/man5/
usr/share/man/man5/syslog-ng.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/syslog-ng.8.gz
usr/share/syslog-ng/
usr/share/syslog-ng/include/
usr/share/syslog-ng/include/scl/
usr/share/syslog-ng/include/scl/pacct/
usr/share/syslog-ng/include/scl/pacct/plugin.conf
usr/share/syslog-ng/include/scl/syslogconf/
usr/share/syslog-ng/include/scl/syslogconf/README
usr/share/syslog-ng/include/scl/syslogconf/convert-syslogconf.awk
usr/share/syslog-ng/include/scl/syslogconf/plugin.conf
usr/share/syslog-ng/include/scl/system/
usr/share/syslog-ng/include/scl/system/generate-system-source.sh
usr/share/syslog-ng/include/scl/system/plugin.conf
usr/share/syslog-ng/tools/
usr/share/syslog-ng/tools/cfg-grammar.y
usr/share/syslog-ng/tools/lex-rules.am
usr/share/syslog-ng/tools/merge-grammar.pl
usr/share/syslog-ng/xsd/
usr/share/syslog-ng/xsd/patterndb-1.xsd
usr/share/syslog-ng/xsd/patterndb-2.xsd
usr/share/syslog-ng/xsd/patterndb-3.xsd
var/
var/lib/
var/lib/syslog-ng/

%BACKUP%
etc/syslog-ng/modules.conf	d341575c0138f5f293bc65744b220d5b
etc/syslog-ng/scl.conf	8f8362660d66f42fdd41a46562fe311a
etc/syslog-ng/syslog-ng.conf	58fb52d3c0b2a6a5a2681571cbd36a44
etc/conf.d/syslog-ng	af23a5f5d39d65c78d7bd898cb0ff9f8
etc/logrotate.d/syslog-ng	0c0d5149bd3e8d666c506351edecbe71

