# partial_test2/var/lib/pacman/local/slim-1.3.3-5/files
%FILES%
etc/
etc/logrotate.d/
etc/logrotate.d/slim
etc/pam.d/
etc/pam.d/slim
etc/rc.d/
etc/rc.d/slim
etc/slim.conf
usr/
usr/bin/
usr/bin/slim
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/slim.service
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/slim.1.gz
usr/share/slim/
usr/share/slim/themes/
usr/share/slim/themes/default/
usr/share/slim/themes/default/background.jpg
usr/share/slim/themes/default/panel.png
usr/share/slim/themes/default/slim.theme

%BACKUP%
etc/slim.conf	ae7554ad60d5c0ee742ed6e0d111d96d
etc/logrotate.d/slim	43da096480bf72c3ccec8ad8400f34f0
etc/pam.d/slim	d33edc74724c6ca00445767ce38fc732

