# ext2fs/var/lib/pacman/local/logrotate-3.8.6-1/files
%FILES%
etc/
etc/cron.daily/
etc/cron.daily/logrotate
etc/logrotate.conf
etc/logrotate.d/
usr/
usr/bin/
usr/bin/logrotate
usr/share/
usr/share/man/
usr/share/man/man5/
usr/share/man/man5/logrotate.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/logrotate.8.gz

%BACKUP%
etc/logrotate.conf	86209d257c8b8bc0ae34d6f6ef057c0f

